# daemon/opencode-gateway.service
[Unit]
Description=TheSmartAI Gateway (C root daemon)
Documentation=https://github.com/sst/opencode
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/opencode-gateway
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=10

# All config from single source of truth
EnvironmentFile=/etc/opencode/opencode.cfg

# Runtime directory for JWT secret, PID tracking
RuntimeDirectory=opencode-gateway
RuntimeDirectoryMode=0755

# Security hardening (gateway needs root for PAM + setuid, but limit the rest)
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/run/opencode-gateway /run/user /tmp
ProtectHome=no

[Install]
WantedBy=multi-user.target
